# analytics.service
[Unit]
Description="Run goaccess to generate report"

[Service]
ExecStart=/usr/bin/goaccess -o /output/path/analytics.html
